# com.supernote.sync.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.supernote.sync.script</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/jeremythake/Code/supernote-to-todoist/src/run_sync.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/tmp/supernote-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/supernote-launchd.log</string>
</dict>
</plist>